# etc/containers/systemd/vaultwarden.container
[Unit]
Description=Vaultwarden
PartOf=vaultwarden.target
After=vaultwarden-postgres.service
Wants=vaultwarden-postgres.service

[Container]
ContainerName=vaultwarden
Image=ghcr.io/polarix-containers/vaultwarden:latest
AutoUpdate=registry

PublishPort=8443:8443/tcp
Network=vaultwarden.network

Volume=/srv/vaultwarden/vaultwarden:/data:Z
Volume=/srv/vaultwarden/certs:/certs:Z,ro
ReadOnly=true

EnvironmentFile=/srv/vaultwarden/vaultwarden.env
User=200003
Group=200003

NoNewPrivileges=true
DropCapability=ALL

[Service]
LockPersonality=true
ProtectClock=true
RestrictRealtime=true
RuntimeDirectoryMode=700
SystemCallArchitectures=native
SystemCallFilter=~@obsolete

Restart=always

[Install]
WantedBy=vaultwarden.target
