# etc/systemd/system/forgejo-sync.service
[Unit]
Description=Synchronize Forgejo certificates
Requisite=forgejo.target

[Service]
Type=oneshot   
ExecStart=/usr/sbin/systemctl restart forgejo