# etc/systemd/system/miniflux-sync.path
[Unit]
Description=Monitor Miniflux certificates

[Path]
PathChanged=/srv/miniflux/certs
Unit=miniflux-sync.service

[Install]
WantedBy=multi-user.target