# etc/systemd/system/opencloud-sync.path
[Unit]
Description=Monitor OpenCloud certificates

[Path]
PathChanged=/srv/opencloud/certs
Unit=opencloud-sync.service

[Install]
WantedBy=multi-user.target