# etc/systemd/system/vaultwarden-sync.path
[Unit]
Description=Monitor Vaultwarden certificates

[Path]
PathChanged=/srv/vaultwarden/certs
Unit=vaultwarden-sync.service

[Install]
WantedBy=multi-user.target